# scripts/deployment/supervisor.conf
# supervisor.conf - Production process management for RV Agentic Workers
#
# Installation:
#   1. Install supervisor: pip install supervisor
#   2. Copy this file: sudo cp supervisor.conf /etc/supervisor/conf.d/rv_agentic.conf
#   3. Update paths below to match your deployment
#   4. Reload supervisor: sudo supervisorctl reread && sudo supervisorctl update
#   5. Start all: sudo supervisorctl start rv_agentic:*
#
# Commands:
#   sudo supervisorctl status rv_agentic:*        # Check status
#   sudo supervisorctl restart rv_agentic:*       # Restart all
#   sudo supervisorctl tail -f rv_agentic:lead_list  # View logs

[group:rv_agentic]
programs=rv_heartbeat_monitor,rv_lead_list_runner,rv_company_research_runner,rv_contact_research_runner

# Heartbeat Monitor - Monitors worker health and cleans up dead workers
[program:rv_heartbeat_monitor]
command=/Users/marklerner/RV_Agentic_FrontEnd_Dev/.venv/bin/python -m rv_agentic.workers.heartbeat_monitor
directory=/Users/marklerner/RV_Agentic_FrontEnd_Dev
user=marklerner
autostart=true
autorestart=true
startretries=10
redirect_stderr=true
stdout_logfile=/Users/marklerner/RV_Agentic_FrontEnd_Dev/logs/supervisor_heartbeat_monitor.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=PATH="/Users/marklerner/RV_Agentic_FrontEnd_Dev/.venv/bin:%(ENV_PATH)s"
stopwaitsecs=30
priority=100

# Lead List Runner - Company discovery worker
[program:rv_lead_list_runner]
command=/Users/marklerner/RV_Agentic_FrontEnd_Dev/.venv/bin/python -m rv_agentic.workers.lead_list_runner
directory=/Users/marklerner/RV_Agentic_FrontEnd_Dev
user=marklerner
autostart=true
autorestart=true
startretries=10
redirect_stderr=true
stdout_logfile=/Users/marklerner/RV_Agentic_FrontEnd_Dev/logs/supervisor_lead_list_runner.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=PATH="/Users/marklerner/RV_Agentic_FrontEnd_Dev/.venv/bin:%(ENV_PATH)s"
stopwaitsecs=60
priority=200

# Company Research Runner - Company enrichment worker
[program:rv_company_research_runner]
command=/Users/marklerner/RV_Agentic_FrontEnd_Dev/.venv/bin/python -m rv_agentic.workers.company_research_runner
directory=/Users/marklerner/RV_Agentic_FrontEnd_Dev
user=marklerner
autostart=true
autorestart=true
startretries=10
redirect_stderr=true
stdout_logfile=/Users/marklerner/RV_Agentic_FrontEnd_Dev/logs/supervisor_company_research_runner.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=PATH="/Users/marklerner/RV_Agentic_FrontEnd_Dev/.venv/bin:%(ENV_PATH)s"
stopwaitsecs=60
priority=300

# Contact Research Runner - Contact discovery worker
[program:rv_contact_research_runner]
command=/Users/marklerner/RV_Agentic_FrontEnd_Dev/.venv/bin/python -m rv_agentic.workers.contact_research_runner
directory=/Users/marklerner/RV_Agentic_FrontEnd_Dev
user=marklerner
autostart=true
autorestart=true
startretries=10
redirect_stderr=true
stdout_logfile=/Users/marklerner/RV_Agentic_FrontEnd_Dev/logs/supervisor_contact_research_runner.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=PATH="/Users/marklerner/RV_Agentic_FrontEnd_Dev/.venv/bin:%(ENV_PATH)s"
stopwaitsecs=60
priority=400

# Optional: Staging Promotion Runner - Uncomment to enable
# [program:rv_staging_promotion_runner]
# command=/Users/marklerner/RV_Agentic_FrontEnd_Dev/.venv/bin/python -m rv_agentic.workers.staging_promotion_runner
# directory=/Users/marklerner/RV_Agentic_FrontEnd_Dev
# user=marklerner
# autostart=true
# autorestart=true
# startretries=10
# redirect_stderr=true
# stdout_logfile=/Users/marklerner/RV_Agentic_FrontEnd_Dev/logs/supervisor_staging_promotion_runner.log
# stdout_logfile_maxbytes=10MB
# stdout_logfile_backups=5
# environment=PATH="/Users/marklerner/RV_Agentic_FrontEnd_Dev/.venv/bin:%(ENV_PATH)s"
# stopwaitsecs=60
# priority=500

# Configuration Notes:
# - autostart=true: Workers start automatically when supervisor starts
# - autorestart=true: Workers restart automatically if they crash
# - startretries=10: Try to restart up to 10 times before giving up
# - stopwaitsecs: How long to wait for graceful shutdown before killing
# - priority: Start order (lower number = starts first)
#
# Supervisor will:
# 1. Start heartbeat monitor first (monitors other workers)
# 2. Start lead list runner (entry point of pipeline)
# 3. Start company research runner
# 4. Start contact research runner
# 5. Automatically restart any worker that crashes
# 6. Rotate logs to prevent disk space issues
